# deploy/gotosocial/rss-sync.service.example
[Unit]
Description=Sync fft.im notes RSS to GoToSocial
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/gotosocial
EnvironmentFile=/opt/gotosocial/rss-sync.env
ExecStart=/usr/bin/node /opt/gotosocial/scripts/sync-rss-to-gotosocial.mjs
