# deploy/gotosocial/rss-sync.timer.example
[Unit]
Description=Run fft.im notes RSS sync every 10 minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=10m
Unit=rss-sync.service
Persistent=true

[Install]
WantedBy=timers.target
